# backend/systemd/arome-ingest.timer
[Unit]
Description=Run AROME-PI / AROME-NWC ingest hourly at HH:10 UTC

[Timer]
# Fire at HH:10:00 — AROME-PI runs at HH:00, files typically appear by HH:05-HH:10.
OnCalendar=*-*-* *:10:00
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
