# backend/systemd/radar-nowcast.timer
[Unit]
Description=Run pysteps nowcast 2 min after each radar-ingest slot

[Timer]
# Fire at HH:02:30, HH:07:30, ... — 2 min after the matching radar-ingest slot
# so radar-ingest has had time to download + render the new observed frame.
OnCalendar=*-*-* *:2/5:30
AccuracySec=15s
Persistent=true

[Install]
WantedBy=timers.target
