# backend/systemd/radar-ingest.service
[Unit]
Description=Météo-France radar mosaic ingestion
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=radar
Group=radar
WorkingDirectory=/opt/radar
EnvironmentFile=/etc/radar/env
ExecStart=/opt/radar/backend/.venv/bin/radar ingest-radar
Nice=10
TimeoutStartSec=120
ProtectSystem=strict
ReadWritePaths=/var/lib/radar
PrivateTmp=true
NoNewPrivileges=true
